# deploy/systemd/units/drift-asyncaiflow.service
[Unit]
Description=Drift AsyncAIFlow Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__STACK_USER__
EnvironmentFile=/etc/drift-stack.env
WorkingDirectory=/opt/drift-stack-systemd
ExecStart=/opt/drift-stack-systemd/bin/run-asyncaiflow-runtime.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target